# autoinstall/first-boot.service
# /etc/systemd/system/jarvis-os-first-boot.service
#
# J.A.R.V.I.S-OS First-Boot Setup Service
# Läuft einmalig nach erstem Proxmox-Boot und deployt die Basis-Plattform.
#
# Installation (wird von build-iso.sh automatisch injiziert):
#   cp first-boot.service /etc/systemd/system/jarvis-os-first-boot.service
#   cp first-boot.sh /root/jarvis-os-first-boot.sh
#   chmod +x /root/jarvis-os-first-boot.sh
#   systemctl enable jarvis-os-first-boot.service
#
# Logs verfolgen:
#   journalctl -u jarvis-os-first-boot -f
#   cat /var/log/jarvis-os-first-boot.log

[Unit]
Description=J.A.R.V.I.S-OS First-Boot Setup (Nginx Proxy + CasaOS + Deployment Hub)
Documentation=https://github.com/WaR10ck-2025/wine-docker-manager

# Nach Netzwerk und Proxmox-Cluster warten
After=network-online.target pve-cluster.service pve-storage.target
Wants=network-online.target

# LUKS muss entsperrt sein (pve-storage.target erfordert das bereits)
Requires=pve-storage.target

# Nur einmal ausführen — wenn Done-Flag existiert, überspringen
ConditionPathExists=!/etc/jarvis-setup.done

[Service]
Type=oneshot
ExecStart=/root/jarvis-os-first-boot.sh
StandardInput=tty
TTYPath=/dev/tty1
TTYReset=yes
StandardOutput=journal+console
StandardError=journal+console
SyslogIdentifier=jarvis-os-first-boot

# Timeout: 30 Minuten (Docker Image Downloads können dauern)
TimeoutStartSec=1800

# Bei Fehler: nicht automatisch neu starten (manuell debuggen)
Restart=no

# Logging auch in Datei
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
